# falconeye.service
[Unit]
Description=FalconEye v3.5.0
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/falconeye/app_src
EnvironmentFile=/opt/falconeye/.env
ExecStart=/opt/falconeye/venv/bin/gunicorn app.main:app \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --worker-class uvicorn.workers.UvicornWorker \
    --timeout 90 \
    --access-logfile /var/log/falconeye/access.log \
    --error-logfile /var/log/falconeye/error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
