# deploy/systemd/falconeye-apnic.service
[Unit]
Description=FalconEye APNIC delegated stats ingest worker
After=network-online.target
Wants=network-online.target
OnSuccess=falconeye-ssg.service

[Service]
Type=oneshot
User=falconeye
Group=falconeye
WorkingDirectory=/opt/falconeye/src
EnvironmentFile=/opt/falconeye/config/secrets.env
ExecStart=/opt/falconeye/venv/bin/python -m falconeye.ingest.apnic
StandardOutput=journal
StandardError=journal
SyslogIdentifier=falconeye-apnic

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/falconeye/db /opt/falconeye/logs
