# deploy/systemd/falconeye-ssg.service
[Unit]
Description=FalconEye static site generator
After=network-online.target

[Service]
Type=oneshot
User=falconeye
Group=falconeye
WorkingDirectory=/opt/falconeye/src
EnvironmentFile=/opt/falconeye/config/secrets.env
ExecStart=/opt/falconeye/venv/bin/python -m falconeye.ssg
StandardOutput=journal
StandardError=journal
SyslogIdentifier=falconeye-ssg

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/falconeye/db /opt/falconeye/public /opt/falconeye/logs
