# deploy/systemd/falconeye-urlhaus.timer
[Unit]
Description=Run FalconEye URLhaus ingest every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
