# deploy/systemd/falconeye-digest.service
[Unit]
Description=FalconEye Ghost daily digest publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=falconeye
Group=falconeye
WorkingDirectory=/opt/falconeye/src
EnvironmentFile=/opt/falconeye/config/secrets.env
ExecStart=/opt/falconeye/venv/bin/python -m falconeye.digest
StandardOutput=journal
StandardError=journal
SyslogIdentifier=falconeye-digest

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/falconeye/db /opt/falconeye/logs
